# global-proxy.service
[Unit]
Description=Global proxy
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/global-proxy init
ExecStop=/usr/local/bin/global-proxy fini
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
